# .github/actions-runner-launchd.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Sample launchd plist for a self-hosted GitHub Actions runner on macOS.
  Generated by `./svc.sh install` (writes to ~/Library/LaunchAgents/).
  This file is the canonical reference for new runners; the actual
  installed plist lives in ~/Library/LaunchAgents/.

  Setup steps (one-time per runner):
    1. mkdir -p ~/actions-runner-vibeflow && cd $_
    2. tar -xzf ~/actions-runner-osx-arm64-2.335.1.tar.gz
    3. Get a registration token from
       https://github.com/magicpro97/vibeflow/settings/actions/runners/new
    4. ./config.sh --unattended \
         --url https://github.com/magicpro97/vibeflow \
         --token <token> \
         --labels "self-hosted,linux,x64" \
         --name "MacBook-Pro-cua-mac-vibeflow"
    5. ./svc.sh install   # writes ~/Library/LaunchAgents/...plist
    6. ./svc.sh start

  The runner now starts automatically on every Mac boot/login.

  To remove: ./svc.sh uninstall
-->
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>actions.runner.magicpro97-vibeflow.MacBook-Pro-cua-mac-vibeflow</string>
    <key>ProgramArguments</key>
    <array>
      <string>/Users/linhn/actions-runner-vibeflow/runsvc.sh</string>
    </array>
    <key>UserName</key>
    <string>linhn</string>
    <key>WorkingDirectory</key>
    <string>/Users/linhn/actions-runner-vibeflow</string>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/linhn/Library/Logs/actions.runner.magicpro97-vibeflow.MacBook-Pro-cua-mac-vibeflow/stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/linhn/Library/Logs/actions.runner.magicpro97-vibeflow.MacBook-Pro-cua-mac-vibeflow/stderr.log</string>
    <key>EnvironmentVariables</key>
    <dict>
      <key>ACTIONS_RUNNER_SVC</key>
      <string>1</string>
    </dict>
    <key>ProcessType</key>
    <string>Interactive</string>
    <key>SessionCreate</key>
    <true/>
  </dict>
</plist>
